# contrib/systemd_fcash.service
[Unit]
Description=BWDB
Requires=network.target

[Service]
Type=simple
WorkingDirectory=/usr/opt/fcash
ExecStart=/usr/bin/bwdb 
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=15
User=fcash
ExecStartPre=/bin/mkdir -p /run/bwdb
ExecStartPre=/bin/chown fcash:fcash /run/bwdb
ExecStartPre=/bin/chmod 755 /run/bwdb
PermissionsStartOnly=true
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
